# vyos-fw-gui.service
[Unit]
Description=VyOS Firewall GUI
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/vyos-fw-gui
ExecStart=/usr/bin/env python3 app.py
Restart=always

[Install]
WantedBy=multi-user.target